# system_files/etc/systemd/system/greetd-useradd.service
[Unit]
Description=Create greetd user on system boot
Before=greetd.service

[Service]
Type=oneshot
ExecStart=/usr/etc/greetd/greetd-useradd
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
